# packaging/scripts/postinstall.sh
systemctl daemon-reload
systemctl enable rds-enhanced-monitoring-exporter.service
systemctl restart rds-enhanced-monitoring-exporter.service
