# packaging/scripts/preremove.sh
systemctl stop rds-enhanced-monitoring-exporter.service
systemctl disable rds-enhanced-monitoring-exporter.service
